# deployment_scripts/puppet/modules/onos/templates/debian/onos.service.erb
[Unit]
Description=onos

[Service]
Type=forking
ExecStart=/etc/init.d/onos start

[Install]
WantedBy=multi-user.target
